# mkosi/base/mkosi.profiles/c8s/mkosi.extra/etc/systemd/system/models-disk.service
# Mount a pre-populated read-only weights disk (serial=confai-models) at
# /var/lib/models before rke2 starts, so large model weights persist across CVM
# relaunches instead of re-downloading. See models-disk.sh for the why.
#
# The disk is public (weights), so unlike containerd-data-disk.service there is
# no encryption and no per-boot mkfs — the filesystem is mounted read-only as-is.
# Absent disk → clean no-op (GPU-less / no-weights boots).
#
# A oneshot (not a .mount unit) because the device is matched by serial at
# runtime, which a static .mount cannot express. RemainAfterExit ties the mount
# lifecycle to the unit. Ordered Before rke2 so /var/lib/models is ready when
# pods start.
[Unit]
Description=Mount the read-only weights disk (serial=confai-models) at /var/lib/models
DefaultDependencies=no
Conflicts=umount.target
Before=rke2-server.service rke2-agent.service
After=systemd-udevd.service
# Order After=systemd-udevd so udevd has started before we resolve /dev/disk/by-id.
# NOT a hard Require, and no udev-settle: a dep on the by-id device unit would hang
# when no disk is attached, breaking the no-op path. The script polls for the serial.

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/models-disk.sh
# Non-fatal to boot: the script wraps the mount in `timeout` and always exits 0
# (a missing/failed disk just means a cold cache). Cap total runtime as a backstop.
TimeoutStartSec=60
ExecStop=-/bin/sh -c 'umount /var/lib/models 2>/dev/null; true'

[Install]
WantedBy=local-fs.target
