# mkosi/base/mkosi.profiles/c8s/mkosi.extra/etc/systemd/system/rke2-server.service.d/no-modprobe.conf
# Upstream rke2-server.service has
#   ExecStartPre=-/sbin/modprobe br_netfilter
#   ExecStartPre=-/sbin/modprobe overlay
# These are intentionally failure-tolerant (the `-` prefix means "don't
# error on non-zero exit"), but on this image they can never succeed and
# only log noise: the only loadable modules are the signed NVIDIA ones, and
# once nvidia-modules-latch.service sets kernel.modules_disabled=1 any
# modprobe attempt fails outright. The features themselves are built-in:
#   - br_netfilter functionality: CONFIG_BRIDGE_NETFILTER=y (kernel/c8s.config)
#   - overlay functionality:      CONFIG_OVERLAY_FS=y       (required.config)
# So we override the ExecStartPre to a no-op. The empty `ExecStartPre=` in
# a drop-in CLEARS the list, then we don't re-add anything.
[Service]
ExecStartPre=
