# mkosi/base/mkosi.profiles/ssh/mkosi.extra/etc/systemd/system/ssh.service.d/no-vsock.conf
# vsock fence: on GPU images the kernel enables AF_VSOCK solely for the
# attestation service's TDX GetQuote path (see steep kernel/gpu.config).
# sshd has no business on the host<->guest channel — deny the family.
# Harmless on non-vsock kernels (the family is simply already absent).
[Service]
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
